# config/launchd/com.compound.review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.compound.review</string>
  
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/vitruvius/git/compound-learning-system/scripts/claude-compound-review.sh</string>
  </array>
  
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>23</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  
  <key>StandardOutPath</key>
  <string>/Users/vitruvius/git/logs/claude-compound-review.log</string>
  
  <key>StandardErrorPath</key>
  <string>/Users/vitruvius/git/logs/claude-compound-review.log</string>
  
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
